# webadmin/files/vnstati.cron
* * * * * root test -x /bin/gen-vnstati-images && /bin/gen-vnstati-images
